# ft_services/srcs/nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = /bin/sh -c "nginx -g 'daemon off;'"

[program:sshd]
command = /bin/sh -c "/usr/sbin/sshd -D"